# templates/postgresql/monitoring/pgwatch2/pgwatch2.service.erb
[Unit]
Description=pgwatch2 gathering daemon
After=network-online.target

[Service]
User=pgwatch2
Type=notify
ExecStart=/usr/bin/pgwatch2-daemon -c /etc/pgwatch2/config/instances.yaml -m /etc/pgwatch2/metrics --datastore=prometheus --prometheus-port=9187 --prometheus-namespace=pgwatch2 --verbose=WARN
Restart=on-failure
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
